# zenbot-config/zenbot-insight/crontab
*/15 * * * * cd path_to_zenbot/zenbot/ && ./zenbot.sh backfill binance.ETH-USDT --days 15

0 0 * * * cd path_to_zenbot/zenbot-insight/ && ./run_strategies.sh
